# scripts/sherpa-asr-multi.service
[Unit]
Description=Sherpa-ONNX Multi-Worker (zh/en INT8)
After=network-online.target

[Service]
Type=simple
User=root
Environment=CUDA_DEVICE_MAX_CONNECTIONS=32
Environment=WORKERS=3
Environment=BASE_PORT=8001
ExecStart=/bin/bash /opt/04_run_server_multi_int8.sh
Restart=always
RestartSec=2
LimitNOFILE=1048576

[Install]
WantedBy=multi-user.target
